# board/rockchip/rv1108/fs-overlay/etc/profile.d/RkEnv.sh
#! /bin/sh
export PATH=/bin:/sbin:/usr/bin:/usr/sbin
export LD_LIBRARY_PATH=/lib/:/usr/lib/:/usr/local/lib/

mkdir /data
mount /dev/mmcblk0p5  /data/


