# package/rockchip/usbhost/61-usb-disk-auto-mount.rules
# Rules for USBmount -*- conf -*-

KERNEL=="sda", SUBSYSTEM=="block", ACTION=="add", RUN+="/bin/ln -s /media/usb0 /udisk"
